# supervisord.conf
[supervisord]
nodaemon=true
user=user
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:chromadb]
# Running Chroma as a server inside the same container
command=python -m chromadb.bin.server --path ./src/chroma_data --port 8000
autostart=true
autorestart=true
stdout_logfile=/app/logs/chroma.out.log
stderr_logfile=/app/logs/chroma.err.log

[program:gradio_app]
# Running your Gradio UI
command=python src/GradioUI.py
autostart=true
autorestart=true
stdout_logfile=/app/logs/app.out.log
stderr_logfile=/app/logs/app.err.log